# .config/systemd/user/home-steph-Owncloud.automount
# Works only as root under /etc/systemd/system/
[Unit]
Description=Automount Owncloud

[Automount]
Where=/home/steph/Owncloud

[Install]
WantedBy=multi-user.target
