# .config/systemd/user/home-steph-Node-Data.automount
# Works only as root under /etc/systemd/system/
[Unit]
Description=Automount Data from Node

[Automount]
Where=/home/steph/Node/Data

[Install]
WantedBy=multi-user.target
